# home/dot_local/share/agentic-workstation/dev-companion/systemd/user/dots-workstation-dev-companion-worker.timer
[Unit]
Description=Timer for agentic-workstation dev companion worker
Requires=dots-workstation-dev-companion-worker.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=dots-workstation-dev-companion-worker.service
Persistent=true

[Install]
WantedBy=timers.target
