# home/dot_local/share/agentic-workstation/dev-companion/systemd/user/dots-workstation-dev-companion-worker.service
[Unit]
Description=agentic-workstation dev companion queue worker (oneshot)
[Service]
Type=oneshot
ExecStart=%h/.local/share/agentic-workstation/dev-companion/worker.sh
